# systemd/system/cron-hourly.target
[Unit]
Description=Hourly cron target
StopWhenUnneeded=true
