# systemd/system/cron-daily.timer
[Unit]
Description=Daily cron script timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=cron-daily.target

[Install]
WantedBy=timers.target
