# systemd/system/cron-daily.target
[Unit]
Description=Daily cron target
StopWhenUnneeded=true
